# roles/manager/files/backup@.service
[Unit]
Description=Backup Job

[Service]
Type=oneshot
WorkingDirectory=/home/%I
ExecStart=/usr/bin/site-backup.sh
User=%I
Group=%I
